# target/etc/supervisor.d/clamonacc.conf
[program:clamonacc]
command=/usr/local/bin/clamonacc.sh
user=root
autostart=true
autorestart=unexpected
stdout_logfile=/var/log/clamav/clamonacc.log
stderr_logfile=/var/log/clamav/clamonacc.err.log
