# target/etc/supervisor.d/entrypoint.conf
[program:docker_entrypoint.sh]
command=/usr/local/bin/docker_entrypoint.sh
user=root
autorestart=unexpected
stdout_logfile=/var/log/docker.log
stdout_logfile_maxbytes=4096
