# target-clamav/etc/supervisor.d/clamd.conf
[program:clamd]
command=/usr/sbin/clamd
user=root
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
