# visualization/supervisor.ini
[program:prometheus]
command=prometheus --config.file="/etc/prometheus/prometheus.yml" --web.listen-address="0.0.0.0:9090"
username=prometheus
priority=1
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:grafana]
command=grafana-server --homepath=/usr/share/grafana --config=/etc/grafana/grafana.ini \
  --packaging=docker cfg:default.log.mode=console \
  cfg:default.paths.data=/var/lib/grafana \
  cfg:default.paths.logs=/var/log/grafana \
  cfg:default.paths.plugins=/var/lib/grafana/plugins \
  cfg:default.paths.provisioning=/etc/grafana/provisioning
username=grafana
priority=2
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true